# deploy/systemd/moli-agent.service
[Unit]
Description=__APP_NAME__ agent service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=__RUN_USER__
Group=__RUN_GROUP__
WorkingDirectory=__INSTALL_PREFIX__
EnvironmentFile=__CONFIG_DIR__/agent.env
ExecStart=__INSTALL_PREFIX__/bin/agent
Restart=always
RestartSec=3
TimeoutStopSec=15
KillMode=control-group
NoNewPrivileges=true
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
